# rootfs/etc/s6-overlay/s6-rc.d/xray-log/run
#!/command/with-contenv sh
exec /command/s6-log -b n0 s10485760 /var/log/xray
